# TarBackupFiles.service.template
[Unit]
Description=Tar Files Backup Service

[Service]
Type=simple
ExecStart=sh %s/TarBackupFiles.sh

